# server/services/asset-report-risks.timer
[Unit]
Description=Run Asset Report Risks Every Hour

[Timer]
OnUnitInactiveSec=1h

[Install]
WantedBy=timers.target
